# cron.d/sync-mingw-packages-from-obs
0 0 * * * <install-root>/tools/sync-mingw-packages-from-obs
